# files/supervisor/syslog.conf
; supervisor configuration for syslogd
[program:syslog]
command=/sbin/syslogd -u syslog -n
